# TEST_OPENSMTPD_DOVECOT/supervisor_dovecot.conf
[program:dovecot]
command=/usr/sbin/dovecot -F -c /etc/dovecot/standalone.conf
environment=IMAP_USER1=%(ENV_IMAP_USER1)s,IMAP_USER2=%(ENV_IMAP_USER2)s,IMAP_PASSWD=%(ENV_IMAP_PASSWD)s
